# packaging/debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create system user and group if they don't exist
        if ! getent group dmx-bridge >/dev/null 2>&1; then
            groupadd --system dmx-bridge
        fi
        if ! id dmx-bridge >/dev/null 2>&1; then
            useradd --system --gid dmx-bridge --home-dir /var/lib/dmx-bridge \
                    --shell /usr/sbin/nologin --comment "DMX LAN Bridge service user" \
                    --no-create-home dmx-bridge
        fi

        # Ensure data directory exists with correct ownership
        if [ ! -d /var/lib/dmx-bridge ]; then
            install -d -o dmx-bridge -g dmx-bridge -m 750 /var/lib/dmx-bridge
        else
            chown dmx-bridge:dmx-bridge /var/lib/dmx-bridge
            chmod 750 /var/lib/dmx-bridge
        fi

        # Ensure config directory has correct permissions
        if [ -d /etc/dmx-bridge ]; then
            chown root:dmx-bridge /etc/dmx-bridge
            chmod 755 /etc/dmx-bridge
        fi

        # Ensure config file has correct permissions (if it exists)
        if [ -f /etc/dmx-bridge/config.toml ]; then
            chown root:dmx-bridge /etc/dmx-bridge/config.toml
            chmod 640 /etc/dmx-bridge/config.toml
        fi

        # Ensure capability catalog directory has correct permissions
        if [ -d /usr/share/dmx_lan_bridge ]; then
            chown root:root /usr/share/dmx_lan_bridge
            chmod 755 /usr/share/dmx_lan_bridge
        fi

        # Ensure capability catalog has correct ownership and permissions
        if [ -f /usr/share/dmx_lan_bridge/capability_catalog_govee.json ]; then
            chown root:root /usr/share/dmx_lan_bridge/capability_catalog_govee.json
            chmod 644 /usr/share/dmx_lan_bridge/capability_catalog_govee.json

            # Remove ACLs if setfacl is available (Ubuntu 24.04 has ACLs that block access)
            if command -v setfacl >/dev/null 2>&1; then
                echo "Removing ACLs from capability catalog..."
                setfacl -b /usr/share/dmx_lan_bridge/capability_catalog_govee.json 2>/dev/null || true
                setfacl -b /usr/share/dmx_lan_bridge 2>/dev/null || true
            else
                echo "Note: setfacl not available, skipping ACL removal"
            fi

            # Verify file is accessible
            echo "Verifying capability catalog is accessible..."
            ls -la /usr/share/dmx_lan_bridge/capability_catalog_govee.json || echo "WARNING: Cannot list file"
            if [ -r /usr/share/dmx_lan_bridge/capability_catalog_govee.json ]; then
                echo "✓ Capability catalog is readable by root"
            else
                echo "WARNING: Capability catalog is not readable by root!"
            fi

            # Test if dmx-bridge user can read the file
            if command -v sudo >/dev/null 2>&1 && id dmx-bridge >/dev/null 2>&1; then
                if sudo -u dmx-bridge test -r /usr/share/dmx_lan_bridge/capability_catalog_govee.json 2>/dev/null; then
                    echo "✓ Capability catalog is readable by dmx-bridge user"
                else
                    echo "WARNING: Capability catalog is NOT readable by dmx-bridge user!"
                    echo "Directory permissions:"
                    ls -lad /usr/share/dmx_lan_bridge/
                fi
            fi
        fi

        # Ensure systemd service file has correct ownership
        for svc_path in /lib/systemd/system/dmx-bridge.service /usr/lib/systemd/system/dmx-bridge.service; do
            if [ -f "$svc_path" ]; then
                chown root:root "$svc_path"
                chmod 644 "$svc_path"
            fi
        done

        # Reload systemd daemon
        systemctl daemon-reload || true

        # Enable service (but don't start yet if this is initial install)
        systemctl enable dmx-bridge.service || true

        # On upgrade, restart the service if it was running
        # On fresh install, start the service
        if systemctl is-active --quiet dmx-bridge.service; then
            echo "Restarting dmx-bridge service..."
            systemctl restart dmx-bridge.service || true
        else
            echo "Starting dmx-bridge service..."
            systemctl start dmx-bridge.service || true
        fi

        echo ""
        echo "DMX LAN Bridge installed successfully!"
        echo ""
        echo "Service status:"
        systemctl status dmx-bridge.service --no-pager -l || true
        echo ""
        echo "Configuration: /etc/dmx-bridge/config.toml"
        echo "Documentation: /usr/share/doc/dmx-lan-bridge/"
        echo ""
        echo "Check service logs: sudo journalctl -u dmx-bridge.service -f"
        echo "Restart service: sudo systemctl restart dmx-bridge.service"
        echo ""
        ;;
esac

#DEBHELPER#
exit 0
