# packaging/systemd/dmx-bridge-user.service
[Unit]
Description=DMX LAN Bridge (user)
Documentation=https://github.com/mccartyp/govee-artnet-lan-bridge
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment=PATH=%h/.local/bin:%h/.local/sbin:/usr/local/bin:/usr/bin
EnvironmentFile=-%h/.config/dmx-bridge/environment
WorkingDirectory=%h/.local/share/dmx-lan-bridge
ExecStart=%h/.local/bin/dmx-lan-bridge --config %h/.config/dmx-bridge/config.toml
Restart=on-failure
RestartSec=2s
# Uncomment when using sd_notify (e.g., via systemd-notify wrapper)
# Type=notify
# NotifyAccess=main

[Install]
WantedBy=default.target
