# packaging/systemd/dmx-bridge.service
[Unit]
Description=DMX LAN Bridge (system)
Documentation=https://github.com/mccartyp//govee-artnet-lan-bridge
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=dmx-bridge
Group=dmx-bridge
WorkingDirectory=/var/lib/dmx-bridge
EnvironmentFile=-/etc/dmx-bridge/environment
# Debug: Verify files are accessible before starting
ExecStartPre=/bin/bash -c 'echo "=== Pre-start file accessibility check (as dmx-bridge user) ===" && ls -la /usr/share/dmx_lan_bridge/capability_catalog_govee.json && test -r /usr/share/dmx_lan_bridge/capability_catalog_govee.json && echo "✓ Capability catalog is readable" || echo "✗ Capability catalog is NOT readable"'
ExecStart=/usr/bin/dmx-lan-bridge --config /etc/dmx-bridge/config.toml
StateDirectory=dmx-bridge
StateDirectoryMode=0750
ConfigurationDirectory=dmx-bridge
NoNewPrivileges=yes
Restart=on-failure
RestartSec=2s
# Uncomment when using sd_notify (e.g., with a wrapper that signals readiness)
# Type=notify
# NotifyAccess=main

[Install]
WantedBy=multi-user.target
