# packaging/systemd/govee-bridge.service
[Unit]
Description=Govee Artnet LAN Bridge (system)
Documentation=https://github.com/mccartyp//govee-artnet-lan-bridge
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=govee-bridge
Group=govee-bridge
WorkingDirectory=/var/lib/govee-bridge
EnvironmentFile=-/etc/govee-bridge/environment
# Debug: Verify files are accessible before starting
ExecStartPre=/bin/bash -c 'echo "=== Pre-start file accessibility check (as govee-bridge user) ===" && ls -la /usr/share/govee_artnet_lan_bridge/capability_catalog.json && test -r /usr/share/govee_artnet_lan_bridge/capability_catalog.json && echo "✓ Capability catalog is readable" || echo "✗ Capability catalog is NOT readable"'
ExecStart=/usr/bin/govee-artnet-bridge --config /etc/govee-bridge/config.toml
StateDirectory=govee-bridge
StateDirectoryMode=0750
ConfigurationDirectory=govee-bridge
NoNewPrivileges=yes
Restart=on-failure
RestartSec=2s
# Uncomment when using sd_notify (e.g., with a wrapper that signals readiness)
# Type=notify
# NotifyAccess=main

[Install]
WantedBy=multi-user.target
