# packaging/systemd/govee-bridge.service
[Unit]
Description=Govee Artnet LAN Bridge (system)
Documentation=https://github.com/mccartyp//govee-artnet-lan-bridge
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=govee-bridge
Group=govee-bridge
WorkingDirectory=/var/lib/govee-bridge
EnvironmentFile=-/etc/govee-bridge/environment
ExecStart=/usr/local/bin/govee-artnet-bridge --config /etc/govee-bridge/config.toml
StateDirectory=govee-bridge
StateDirectoryMode=0750
ConfigurationDirectory=govee-bridge
NoNewPrivileges=yes
Restart=on-failure
RestartSec=2s
# Uncomment when using sd_notify (e.g., with a wrapper that signals readiness)
# Type=notify
# NotifyAccess=main

[Install]
WantedBy=multi-user.target
